# scripts/github-webhook-deploy.node.service.example
# 若坚持用宿主机 Node 跑 webhook（需 Node 18+），可参考此 unit（勿与 Docker 方案同时启用）。
# 复制为 /etc/systemd/system/github-webhook-node.service 并修改路径。

[Unit]
Description=GitHub webhook (node, repo root)
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/yaoyaoweiba-site
EnvironmentFile=-/home/deploy/yaoyaoweiba-site/deploy/webhook/.env.webhook
ExecStart=/usr/bin/node /home/deploy/yaoyaoweiba-site/scripts/github-webhook-deploy.mjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
