# scripts/github-webhook-deploy.service
# 将本文件复制到 /etc/systemd/system/github-webhook-deploy.service 并修改路径后：
#   sudo systemctl daemon-reload && sudo systemctl enable --now github-webhook-deploy
#
# 使用 deploy/webhook 下的独立 compose；首次需先执行：
#   cd /path/to/repo/deploy/webhook && docker compose build

[Unit]
Description=GitHub webhook (docker compose, deploy/webhook)
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=simple
Restart=always
RestartSec=10
User=root
WorkingDirectory=/home/deploy/yaoyaoweiba-site/deploy/webhook
ExecStart=/usr/bin/docker compose up github-webhook
ExecStop=/usr/bin/docker compose stop -t 5 github-webhook

[Install]
WantedBy=multi-user.target
